# infra/checkup-box/checkup-backup.service
[Unit]
Description=Backup diario do Postgres do checkup para S3
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/checkup/src/infra/checkup-box/backup-pgdump.sh
